# controller/systemd/volteria-watchdog.service
[Unit]
Description=Volteria Hardware Watchdog Feeder
After=volteria-system.service

[Service]
Type=simple
ExecStart=/opt/volteria/scripts/watchdog-feeder.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
