# controller/systemd/volteria-logging.service
[Unit]
Description=Volteria Logging Service (Local DB, Cloud Sync, Alarms)
Documentation=https://volteria.org/docs
After=volteria-control.service
Requires=volteria-control.service
PartOf=volteria-supervisor.service

[Service]
Type=simple
User=volteria
Group=volteria
WorkingDirectory=/opt/volteria/controller
ExecStart=/opt/volteria/venv/bin/python -m services.logging.service
Restart=on-failure
RestartSec=10

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=VOLTERIA_STATE_DIR=/run/volteria/state
EnvironmentFile=/etc/volteria/env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=volteria-logging

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/volteria/data /var/log/volteria /run/volteria
PrivateTmp=yes

[Install]
WantedBy=volteria-supervisor.service
