# controller/systemd/volteria-config.service
[Unit]
Description=Volteria Config Service (Sync, Cache, Versioning)
Documentation=https://volteria.org/docs
After=volteria-system.service
Requires=volteria-system.service
PartOf=volteria-supervisor.service

[Service]
Type=simple
User=volteria
Group=volteria
WorkingDirectory=/opt/volteria/controller
ExecStart=/opt/volteria/venv/bin/python -m services.config.service
Restart=on-failure
RestartSec=10
WatchdogSec=120

# Environment
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/etc/volteria/env

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=volteria-config

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/opt/volteria/data /etc/volteria /var/log/volteria
PrivateTmp=yes

[Install]
WantedBy=volteria-supervisor.service
